# recipes-extended/browser-poc/browser-poc-git/browser.service
[Unit]
Description=Genivi Browser PoC service
Requires=dbus.service

[Service]
Environment=XDG_RUNTIME_DIR=/tmp
ExecStart=/opt/browser/bin/browser
