# recipes-graphics/wayland/wayland-ivi-extension/EGLWLInputEventExample.service
[Unit]
Description=EGL Wayland Input Event Example

[Service]
Environment=XDG_RUNTIME_DIR=/tmp
ExecStart=/usr/bin/EGLWLInputEventExample
